# vep1445-traffic-gen.service
[Unit]
Description=VEP1445 Multi-LAN Traffic Generator and Network Impairment System
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/vep1445-traffic-gen
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 /opt/vep1445-traffic-gen/web_api.py
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/vep1445/traffic-gen.log
StandardError=append:/var/log/vep1445/traffic-gen-error.log

# Security settings
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log/vep1445 /home/claude /opt/vep1445-traffic-gen

# Network capabilities
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
